# plugins/infolead-claude-subscription-router/systemd/claude-overnight-executor.service
[Unit]
Description=Claude Code Overnight Work Executor
Documentation=https://github.com/yannickloth/claude-router-system
After=network.target

[Service]
Type=oneshot
# NOTE: This is a TEMPLATE file
# The setup-overnight-execution.sh script will replace __PLUGIN_ROOT__ with actual path
ExecStart=__PLUGIN_ROOT__/scripts/overnight-executor.sh
StandardOutput=append:%h/.claude/infolead-claude-subscription-router/logs/overnight-executor.log
StandardError=append:%h/.claude/infolead-claude-subscription-router/logs/overnight-executor-error.log
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=__PLUGIN_ROOT__/implementation"
WorkingDirectory=__PLUGIN_ROOT__

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.claude

# Resource limits (prevent runaway execution)
TimeoutStartSec=3h
MemoryMax=2G
CPUQuota=80%

[Install]
WantedBy=default.target
